# system/rm7days.timer
[Unit]
Description=timer rm7days

[Timer]
OnUnitActiveSec=86400

[Install]
WantedBy=multi-user.target